# ops/launchd/com.lumenary.research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lumenary.research</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/johnforrester/spirituality/scripts/run_lumenary_parallel_half_hourly.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johnforrester/spirituality</string>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>/Users/johnforrester/spirituality/runs/parallel-research.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnforrester/spirituality/runs/parallel-research.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/johnforrester/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
