# ops/launchd/com.lumenary.map-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lumenary.map-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/johnforrester/spirituality/scripts/update_living_map.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johnforrester/spirituality</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>/Users/johnforrester/spirituality/runs/map-refresh.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnforrester/spirituality/runs/map-refresh.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LUMENARY_MAP_SCHEDULED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
