# ops/launchd/com.lumenary.daily-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lumenary.daily-publish</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/johnforrester/spirituality/scripts/generate_journal_update.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johnforrester/spirituality</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>45</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/johnforrester/spirituality/runs/daily-publish.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnforrester/spirituality/runs/daily-publish.stderr.log</string>
</dict>
</plist>
